# christmas-community-addon/rootfs/etc/services.d/christmas-community/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Finish script for Christmas Community service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

bashio::log.info "Christmas Community service stopped"
